# deploy/polymarketbot.service
[Unit]
Description=PolymarketBot Copy Trading Dashboard
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/polymarketbot
EnvironmentFile=/root/polymarketbot/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/polymarketbot/venv/bin/python -u run_ghost_mode.py --port 8765
Restart=always
RestartSec=30
StartLimitInterval=600
StartLimitBurst=10
StandardOutput=append:/var/log/polymarketbot.log
StandardError=append:/var/log/polymarketbot.log

[Install]
WantedBy=multi-user.target
